# systemd/hl_runtime.slice
[Unit]
Description=Hyperliquid runtime slice

[Slice]
CPUAccounting=true
MemoryAccounting=true
IOAccounting=true
